# images/blog/2019-06-12-restic-backups/btrfs-scrub@.timer
[Unit]
Description=Monthly Btrfs scrub on %f

[Timer]
OnCalendar=monthly
RandomizedDelaySec=10min
Persistent=true

[Install]
WantedBy=multi-user.target