# images/blog/2019-06-12-restic-backups/restic-weekly@.timer
[Unit]
Description=Restic Backup Timer
 
[Timer]
WakeSystem=false
OnCalendar=weekly
Persistent=true
 
[Install]
WantedBy=timers.target