# images/blog/2019-06-12-restic-backups/restic-weekly@.service
[Unit]
Description=Restic Backup
 
[Service]
Type=simple
Nice=19
IOSchedulingClass=2
IOSchedulingPriority=7
ExecStart=%f/restic/restic-backup.sh